# resources/supervisor/bulksync-worker.conf
; ─────────────────────────────────────────────────────────
; Supervisor config for BulkSync queue workers
; Copy to: /etc/supervisor/conf.d/bulksync-worker.conf
; Then run: supervisorctl reread && supervisorctl update
; ─────────────────────────────────────────────────────────

[program:bulksync-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /opt/homebrew/var/www/bulk/artisan queue:work database --queue=bulk-upload --sleep=1 --tries=3 --timeout=180 --max-jobs=500 --max-time=3600
directory=/opt/homebrew/var/www/bulk
autostart=true
autorestart=true
; ── Number of parallel workers ──────────────────────────
; 8 workers × ~3s/image = ~2.7 images/second
; 20,000 images / 2.7 = ~7,400 seconds (~2 hours)
; Adjust numprocs to match your server CPU cores
numprocs=8
; ─────────────────────────────────────────────────────────
redirect_stderr=true
stdout_logfile=/opt/homebrew/var/www/bulk/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=180
